# launchd/com.consciousness.weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.consciousness.weekly</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>/Users/noelgradisar/.openclaw/workspace/dispenza-prompts/scripts/generate-weekly-summary.js</string>
    </array>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>TELEGRAM_BOT_TOKEN</key>
        <string>8142095199:AAGqE3GlP4xsWDpuGf_iesHyQTjzrzEmqX4</string>
        <key>TELEGRAM_CHAT_ID</key>
        <string>1243017061</string>
    </dict>
    
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>20</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/noelgradisar/.openclaw/workspace/dispenza-prompts/logs/weekly.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/noelgradisar/.openclaw/workspace/dispenza-prompts/logs/weekly.error.log</string>
    
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
